# packaging/hyprshell.service
[Unit]
Description=hyprshell daemon
PartOf=graphical-session.target
After=graphical-session.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
ExecStartPre=/bin/sh -c '[ "$XDG_CURRENT_DESKTOP" = "Hyprland" ] || exit 0'
ExecStart=/usr/bin/hyprshell run
Restart=on-failure
RestartSec=20

[Install]
WantedBy=graphical-session.target
